# script/mew.service
[Unit]
Description=mewsichip
After=default.target

[Service]
ExecStart=/home/chip/mewsichip/start
Type=simple

[Install]
WantedBy=default.target
